# docker/overlay/etc/supervisord.d/nginx.ini
[program:nginx]
autorestart=true
autostart=true
command=/usr/sbin/nginx -g "daemon off;"
directory=/data
process_name=%(program_name)s
stderr_logfile=/data/var/log/%(program_name)s_error.log
stderr_logfile_backups=100
stderr_logfile_maxbytes=10MB
stdout_logfile=/data/var/log/%(program_name)s.log
stderr_logfile_backups=100
stdout_logfile_maxbytes=10MB