# docker/overlay/etc/supervisord.d/sha_apid.ini
[program:sha_apid]
autorestart=true
autostart=true
command=/usr/bin/sha_apid
directory=/data
environment=SHA_API_CONFIG="/data/etc/sha_api/sha_api.conf"
process_name=%(program_name)s
stderr_logfile=/data/var/log/%(program_name)s_error.log
stderr_logfile_backups=100
stderr_logfile_maxbytes=10MB
stdout_logfile=/data/var/log/%(program_name)s.log
stderr_logfile_backups=100
stdout_logfile_maxbytes=10MB